# scheduler.cron
*/30 * * * * python main.py >> /var/log/pad.log 2>&1
# */30 * * * * docker start pad